# etc/systemd/rswitch.service
[Unit]
Description=rSwitch XDP-based Software Switch
After=network.target network-online.target
Wants=network-online.target
Documentation=https://github.com/rswitch
OnFailure=rswitch-failsafe.service
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=notify
NotifyAccess=all
ExecStartPre=-/opt/rswitch/scripts/rswitch-failsafe.sh teardown
ExecStartPre=/opt/rswitch/scripts/rswitch-init.sh prepare
ExecStart=/opt/rswitch/scripts/rswitch-init.sh start
ExecStop=/opt/rswitch/scripts/rswitch-init.sh stop
ExecReload=/opt/rswitch/scripts/rswitch-init.sh reload
Restart=on-failure
RestartSec=5
LimitMEMLOCK=infinity
LimitNOFILE=65535
Environment=RSWITCH_HOME=/opt/rswitch
Environment=RSWITCH_PROFILE=default.yaml
Environment=RSWITCH_INTERFACES=ens34,ens35,ens36,ens37

[Install]
WantedBy=multi-user.target
